# assets/systemd/sanalpanel.service
[Unit]
Description=SanalPanel backend
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sanalpanel/env
ExecStart=/opt/sanalpanel/bin/sanalpanel-server
Restart=on-failure
RestartSec=3
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sanalpanel

[Install]
WantedBy=multi-user.target
