# assets/systemd/sanalpanel-db-backup.service
[Unit]
Description=SanalPanel panel DB yedegi

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sanalpanel-db-backup
